# coredump-daemon.service
[Unit]
Description=CoreDump Activity Tracker Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/coredump-daemon
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment="DISPLAY=:1"
Environment="XAUTHORITY=/run/user/1000/gdm/Xauthority"

[Install]
WantedBy=default.target
